# build_files/tilefin-hypridle-rearm.timer
[Unit]
Description=Re-arm hypridle idle detection at the business-hours boundary (Tilefin S26)

[Timer]
# Fire at the end of business hours on weekdays so a machine that went
# idle during the workday re-evaluates and can auto-suspend (R26.4).
OnCalendar=Mon..Fri 18:00
# A missed firing (machine off at 18:00) must not run on next boot — the
# re-arm is only meaningful at the boundary itself.
Persistent=false

[Install]
WantedBy=timers.target
